# build/etc/service/varnish/run
#!/bin/bash

varnishd -a :80 -F -f /etc/varnish/default.vcl
